# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Wireless starter and Interface checker by Kenvix <i@kenvix.com>

# Start mtk wifi

sleep 5s
# 2.4GHz
ip link set ra0 up
brctl addif br-lan ra0

# 5GHz
ip link set rai0 up
brctl addif br-lan rai0

echo "Wireless interface started"

sleep 2s
lanCheck=`uci get network.lan.ifname`
uci set network.lan.ifname="$lanCheck rai0 ra0"
uci commit


exit 0
